# Daemon/ydotool.service.in
[Unit]
Description=Starts ydotoold service

[Service]
Type=simple
Restart=always
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/ydotoold
ExecReload=/usr/bin/kill -HUP $MAINPID
KillMode=process
TimeoutSec=180

[Install]
WantedBy=default.target
